# templates/default/services/apt-mirror.conf.erb
#!upstart
description "Update apt mirror"

start on started apt-mirror-shim
stop on stopping apt-mirror-shim

task
respawn

setuid <%= node.apt_mirror.user %>

env USER=<%= node.apt_mirror.user %>
env HOME=/home/<%= node.apt_mirror.user %>

exec apt-mirror
